# root/etc/services.d/radicale-daemon/run
#!/command/with-contenv sh

while true; do
  cd /var/radicale
  python3 -m radicale --config=/var/config.ini
  sleep 5m
done
